# systemd/ollama.service
[Unit]
Description=Ollama Service
Documentation=https://ollama.ai
After=network.target
Wants=network.target

[Service]
Type=notify
User=ollama
Group=ollama

# Environment configuration
Environment="OLLAMA_HOST=0.0.0.0:11434"
Environment="OLLAMA_MODELS=/var/lib/ollama/models"

# Working directory
WorkingDirectory=/var/lib/ollama

# Main execution
ExecStart=/usr/local/bin/ollama serve

# Restart configuration
Restart=on-failure
RestartSec=5

# Resource limits
LimitNOFILE=65535
LimitNPROC=unlimited

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

# Allow write access to Ollama directories
ReadWritePaths=/var/lib/ollama

# GPU access (AMD ROCm)
# Uncomment the following lines for AMD GPU support
# DevicePolicy=closed
# DeviceAllow=/dev/kfd rw
# DeviceAllow=/dev/dri rw

# GPU access (NVIDIA CUDA)
# Uncomment the following lines for NVIDIA GPU support
# DevicePolicy=closed
# DeviceAllow=/dev/nvidia0 rw
# DeviceAllow=/dev/nvidiactl rw
# DeviceAllow=/dev/nvidia-uvm rw

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama

[Install]
WantedBy=multi-user.target
